# web/backend/service/rapid-web.service
[Unit]
Description=RAPID Web Dashboard Service
After=network.target
Documentation=https://github.com/fishjojo1/RAPID

[Service]
Type=simple
WorkingDirectory=__RAPID_ROOT__/web/backend
ExecStart=__RAPID_ROOT__/web/backend/.venv/bin/rapid-web
Environment=RAPID_WEB=true
Environment=RAPID_WEB_PORT=8998
Environment=RAPID_WEB_HOST=127.0.0.1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
